# etc/service/influxdb/run
#!/bin/bash -e

exec /usr/bin/influxd -config=/etc/influxdb/influxdb.toml 2>&1
